# root/etc/s6-overlay/s6-rc.d/init-grav-config/up
/etc/s6-overlay/s6-rc.d/init-grav-config/run
